# examples/startup/ssh-cert-authority.conf
# this is for upstart
start on runlevel [2345]
setuid ubuntu
setgid ubuntu

script
    export HOME=/home/ubuntu
    ssh-agent /usr/local/bin/ssh-cert-authority runserver
end script
